# web/deploy/telefone-portal.service
[Unit]
Description=Telefone Fixo - Portal Web para Pais
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/telefone-portal
Environment=TELEFONE_DB_PATH=/opt/telefone-portal/data/telefone.db
Environment=TELEFONE_SECRET_KEY=TROCAR_POR_CHAVE_SECRETA
Environment=TELEFONE_SMTP_USER=rede.telefonefixo@gmail.com
Environment="TELEFONE_SMTP_PASSWORD=kpjy dlya aoml scoi"
ExecStart=/opt/telefone-portal/venv/bin/python -m flask run --host=127.0.0.1 --port=5000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
